# config/includes.chroot/etc/profile.d/debianmoss.sh
export DEBIANMOSS=1
export MOSS_ASSETS_DIR=/usr/share/debianmoss
export MOSS_WALLPAPER=/usr/share/backgrounds/debianmoss/tiles.png
export MOSS_ASCII_ART=/usr/share/debianmoss/ascii_art.json
